# packaging/linux/nullgate-update.service
[Unit]
Description=Nullgate auto-update check
Documentation=https://github.com/steeb-k/nullgate
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/nullgatectl --update
