# packaging/linux/nullgate-daemon.service
[Unit]
Description=Nullgate daemon
Documentation=https://github.com/steeb-k/iroh-private-network
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Runs as root so it can create the TUN interface. AmbientCapabilities is
# belt-and-suspenders (and lets you switch to a non-root User= later).
ExecStart=/usr/local/bin/nullgate-daemon run
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
