# packaging/macos/io.github.steeb_k.Nullgate.gui.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Nullgate tray agent autostart LaunchAgent (per-user, Aqua). Installed to
     /Library/LaunchAgents so it loads for each user at login. Launches the
     lightweight headless tray agent (`--agent`), which owns the tray icon +
     notifications and opens the GUI on demand — not the full GUI. No KeepAlive —
     if the user quits it (tray "Quit Nullgate"), it stays quit until next login.
     The agent runs unprivileged and talks to the root daemon over
     /tmp/nullgate.sock. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.steeb_k.Nullgate.gui</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Nullgate.app/Contents/MacOS/nullgate</string>
    <string>--agent</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
</dict>
</plist>
