# packaging/macos/io.github.steeb_k.IPN.update.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- IPN auto-update LaunchDaemon (SYSTEM/root). Installed to /Library/LaunchDaemons.
     Runs as root because applying an update replaces the .app and restarts the
     daemon. RunAtLoad checks shortly after boot; StartCalendarInterval re-checks
     daily. `ipnctl --update` no-ops fast when already current. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.steeb_k.IPN.update</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/ipnctl</string>
    <string>--update</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>13</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>ProcessType</key>
  <string>Background</string>
  <key>StandardOutPath</key>
  <string>/Library/Logs/IPN/ipn-update.log</string>
  <key>StandardErrorPath</key>
  <string>/Library/Logs/IPN/ipn-update.log</string>
</dict>
</plist>
